# report-new-files.cron
BAKDATE=`date +%Y%m%d`
SERVER=`hostname -f`
40 0 * * * root bash /usr/local/bin/summarize.sh 2>&1
